# launchd/ai.openclaw.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- OpenClaw Health Watchdog — runs every 5 minutes -->
  <!-- Install: cp ai.openclaw.watchdog.plist ~/Library/LaunchAgents/ -->
  <!--          launchctl load ~/Library/LaunchAgents/ai.openclaw.watchdog.plist -->

  <key>Label</key>
  <string>ai.openclaw.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/usr/local/share/openclaw/openclaw-watchdog.sh</string>
  </array>

  <!-- Run every 300 seconds (5 minutes) -->
  <key>StartInterval</key>
  <integer>300</integer>

  <!-- Run immediately on load -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Environment: inherit PATH from Homebrew -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <!-- Uncomment and set for Discord alerts: -->
    <!-- <key>OPENCLAW_WATCHDOG_WEBHOOK</key>
    <string>https://discord.com/api/webhooks/YOUR_WEBHOOK_HERE</string> -->
    <!-- Or use bot token: -->
    <!-- <key>OPENCLAW_DISCORD_BOT_TOKEN</key>
    <string>YOUR_BOT_TOKEN</string>
    <key>OPENCLAW_ALERT_CHANNEL</key>
    <string>YOUR_CHANNEL_ID</string> -->
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/openclaw-watchdog.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/openclaw-watchdog.log</string>

  <!-- Don't thrash on repeated failures -->
  <key>ThrottleInterval</key>
  <integer>60</integer>
</dict>
</plist>
